# packaging/debian/spectrum-knx.service
[Unit]
Description=Spectrum KNX bus monitor and analyzer
Documentation=https://github.com/martinhoefling/SpectrumKNX
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=spectrum-knx
Group=spectrum-knx
# Defaults; overridden by the EnvironmentFile below. @VERSION@ is replaced
# at package build time.
Environment=APP_VERSION=@VERSION@
Environment=BIND_HOST=0.0.0.0
Environment=BIND_PORT=8000
EnvironmentFile=-/etc/spectrum-knx/spectrum-knx.env
WorkingDirectory=/opt/spectrum-knx/app
ExecStart=/opt/spectrum-knx/venv/bin/uvicorn main:app --host ${BIND_HOST} --port ${BIND_PORT}
Restart=on-failure
RestartSec=5

# State (sqlite database, uploaded .knxproj) lives in /var/lib/spectrum-knx
StateDirectory=spectrum-knx
StateDirectoryMode=0750

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/spectrum-knx

[Install]
WantedBy=multi-user.target
